# docker/systemd/webapp.service
[Unit]
Description=WebApp Next.js Application
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/vector/projects/webapp-conversation
ExecStart=/usr/bin/npm run dev
Environment="NODE_ENV=production"
Environment="PORT=3000"
Restart=always
RestartSec=15
User=vector
Group=vector

# Resource limits
MemoryMax=4G
MemoryHigh=3G

[Install]
WantedBy=multi-user.target